# debian/nova-api-metadata.upstart.in
description "Nova Metadata API server"
author "Joe Heck <heckj@mac.com>"

start on (filesystem and net-device-up IFACE!=lo)
stop on runlevel [016]


chdir /var/run

pre-start script
	mkdir -p /var/run/nova
	chown nova:root /var/run/nova/

	mkdir -p /var/lock/nova
	chown nova:root /var/lock/nova/
end script

script
    set +e
    runcon -t nova-api_t /bin/true 2>&1 >/dev/null
    if [ $? -eq 0 ]; then
        exec su -s /bin/sh -c "exec runcon -t nova-api_t nova-api-metadata --config-file=/etc/nova/nova.conf" nova
    else
        exec su -s /bin/sh -c "exec nova-api-metadata --config-file=/etc/nova/nova.conf" nova
    fi
end script
